# backend/manga-rider-backend.service
[Unit]
Description=Manga Rider backend systemd service
DefaultDependencies=no
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/java -jar /home/jenkins/demo.jar
EnvironmentFile=/home/jenkins/backend.conf
TimeoutStartSec=0
RemainAfterExit=yes
Restart=on-failure

[Install]
WantedBy=default.target